# linux-pstore.crontab
@reboot		root	/usr/sbin/linux-pstore
@daily		root	/usr/sbin/linux-pstore

